# other/first-boot.service
[Unit]
Description=First boot

[Service]
ExecStart=/first-boot.bash

[Install]
WantedBy=multi-user.target
